# src/slurm-master/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:munged]
command=/usr/sbin/munged --foreground
user=munge
autostart=true
autorestart=true
stdout_logfile=/var/log/munge/munged.log
stderr_logfile=/var/log/munge/munged.log
priority=100

[program:slurmdbd]
command=/usr/sbin/slurmdbd -D
user=slurm
autostart=true
autorestart=true
stdout_logfile=/var/log/slurm/slurmdbd.log
stderr_logfile=/var/log/slurm/slurmdbd.log
priority=200

[program:slurmctld]
command=/usr/sbin/slurmctld -D
user=slurm
autostart=true
autorestart=true
stdout_logfile=/var/log/slurm/slurmctld.log
stderr_logfile=/var/log/slurm/slurmctld.log
priority=300