# src/slurm-master/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:munged]
command=/usr/sbin/munged --foreground
user=munge
autostart=true
autorestart=true
stdout_logfile=/var/log/munge/munged.log
stderr_logfile=/var/log/munge/munged.log
priority=100

[program:slurmdbd]
command=/bin/bash -c 'if [ -f /opt/slurm-installed ] && [ -f /opt/slurmdbd-path ]; then SLURMDBD_PATH=$(cat /opt/slurmdbd-path); if [ -x "$SLURMDBD_PATH" ]; then echo "启动 slurmdbd: $SLURMDBD_PATH"; exec "$SLURMDBD_PATH" -D; else echo "slurmdbd二进制文件不存在: $SLURMDBD_PATH"; sleep infinity; fi; else echo "SLURM未安装或路径文件缺失，跳过slurmdbd启动"; sleep infinity; fi'
user=slurm
autostart=true
autorestart=true
stdout_logfile=/var/log/slurm/slurmdbd.log
stderr_logfile=/var/log/slurm/slurmdbd.log
priority=200

[program:slurmctld]
command=/bin/bash -c 'if [ -f /opt/slurm-installed ] && [ -f /opt/slurmctld-path ]; then SLURMCTLD_PATH=$(cat /opt/slurmctld-path); if [ -x "$SLURMCTLD_PATH" ]; then echo "启动 slurmctld: $SLURMCTLD_PATH"; exec "$SLURMCTLD_PATH" -D; else echo "slurmctld二进制文件不存在: $SLURMCTLD_PATH"; sleep infinity; fi; else echo "SLURM未安装或路径文件缺失，跳过slurmctld启动"; sleep infinity; fi'
user=slurm
autostart=true
autorestart=true
stdout_logfile=/var/log/slurm/slurmctld.log
stderr_logfile=/var/log/slurm/slurmctld.log
priority=300